# harness/launchd/com.example.harness.weekly.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.harness.weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>REPO_ROOT/harness/scripts/local_update.py</string>
    <string>--mode</string>
    <string>weekly</string>
    <string>--calendar</string>
    <string>--llm</string>
  </array>
  <key>WorkingDirectory</key>
  <string>REPO_ROOT</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>REPO_ROOT/harness/updates/weekly.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>REPO_ROOT/harness/updates/weekly.launchd.err</string>
</dict>
</plist>
